# deploy/spotme-whoop-sync.service
[Unit]
Description=Whoop biometric sync job

[Service]
Type=oneshot
WorkingDirectory=%h/spotme
ExecStart=%h/spotme/.venv/bin/python -m server.cli whoop_sync
